# stalwart/systemd/email-glue.service
[Unit]
Description=Email Glue
After=network-online.target stalwart.service
Wants=network-online.target

[Service]
Type=simple
User=emailglue
Group=emailglue
EnvironmentFile=/etc/sysconfig/email-glue
ExecStart=/usr/local/bin/email-glue
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/email-glue
ReadOnlyPaths=/var/lib/stalwart/certs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
